# charger.timer
[Unit]
Description=Timer to delay charger service start

[Timer]
OnBootSec=60

[Install]
WantedBy=timers.target
